# roles/smart_tests/templates/smart-selftest-short.service.j2
[Unit]
Description=Run SMART short self-test ({{ item.tag | upper }})

[Service]
Type=oneshot
Environment=DEVICE=/dev/disk/by-id/{{ item.device }}
ExecStart=/usr/sbin/smartctl -t short ${DEVICE}
